# configs/supervisor/default.conf
[program:gunicorn]
command=/home/USERNAME/PROJECTNAME/bin/start_gunicorn.sh
user=USERNAME
process_name=%(program_name)s
numproc=1
autostart=1
autorestart=1
redirect_stderr=true
